# root/etc/s6/aria2/run
#!/bin/bash

[[ -f ./setup ]] && source ./setup

exec su-exec "${USERNAME}" aria2c --conf-path=/conf/aria2.conf