# root/etc/s6/tracker/run
#!/bin/bash

pushd /conf > /dev/null
    exec su-exec ${USERNAME} python3 /etc/aria2/shell/aria2.py -c "/conf"
popd